# s6/spawn-fcgi_mdclean/finish
#!/bin/sh
echo "mdclean terminated"
s6-svscanctl -t /s6

