# s6/nginx/finish
#!/bin/sh
echo "nginx terminated"
s6-svscanctl -t /s6

